# src/wasm/templates/systemd/docker-compose.service.j2
# Systemd service for Docker Compose: {{ name }}
# Generated by WASM

[Unit]
Description={{ description }}
Documentation=https://github.com/Perkybeet/wasm
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
WorkingDirectory={{ working_directory }}

# Environment
{% for key, value in environment.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% if compose_file %}
Environment="COMPOSE_FILE={{ compose_file }}"
{% endif %}

# Commands
ExecStartPre=/usr/bin/docker compose pull --ignore-pull-failures
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --remove-orphans --build

# Timeouts
TimeoutStartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ name }}

[Install]
WantedBy=multi-user.target
